# .config/systemd/user/xmonad-session.target

[Unit]
Description=User systemd services for XMonad
Wants=pendingUpdates.service
BindsTo=graphical-session.target
